# root/etc/supervisor.d/crond.conf
[program:crond]
command=/usr/sbin/crond -f -L /dev/stdout
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
redirect_stderr=true
