# root/etc/supervisor.d/trojan-go-server.conf
[program:trojan-go-server]
command=/usr/bin/trojan-go -config /etc/trojan-go/config.yaml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
redirect_stderr=true
